# supervisord.conf
[supervisord]
nodaemon=true

[program:wrapdocker]
environment=PORT=4243
command=wrapdocker

[program:gitlab-ci-runner]
command=/app/init app:start
