# supervisord/wrapdocker.conf
[program:wrapdocker]
environment=PORT=4243
command=wrapdocker
